# sabre_systemd/files/logstash.service
[Unit]
Description=logstash
After=network.target

[Service]
ExecStart=/usr/bin/logstash

[Install]
WantedBy=multi-user.target
